# deploy/vps/systemd/hermes-telegram-export-cleanup.timer
[Unit]
Description=Remove expired JarHert Telegram text exports

[Timer]
OnCalendar=*-*-* 03:35:00
Persistent=true
Unit=hermes-telegram-export-cleanup.service

[Install]
WantedBy=timers.target
