# deploy/vps/systemd/hermes-memory-consolidation.timer
[Unit]
Description=Run JarHert memory consolidation nightly

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
Unit=hermes-memory-consolidation.service

[Install]
WantedBy=timers.target
