# deploy/vps/systemd/hermes-watchdog.service
[Unit]
Description=JarHert Hermes gateway watchdog
After=network-online.target hermes-gateway-jarhert.service
ConditionPathExists=!%h/.hermes/profiles/jarhert/state/maintenance

[Service]
Type=oneshot
WorkingDirectory=%h/.hermes/profiles/jarhert
ExecStart=%h/.hermes/profiles/jarhert/.venv/bin/python %h/.hermes/profiles/jarhert/scripts/watchdog.py --min-disk-free-gib 5 --max-memory-percent 90 --restart-inactive
