# deploy/vps/systemd/hermes-telegram-export-cleanup.service
[Unit]
Description=JarHert Telegram text export retention cleanup

[Service]
Type=oneshot
WorkingDirectory=%h/.hermes/profiles/jarhert
Environment=HERMES_HOME=%h/.hermes/profiles/jarhert
EnvironmentFile=%h/.hermes/profiles/jarhert/.env
ExecStart=%h/.hermes/profiles/jarhert/.venv/bin/python %h/.hermes/profiles/jarhert/scripts/cleanup_telegram_exports.py
